# supervisor/conf/chef-client.ini
[program:chef-client]
command=bash -c "/usr/bin/chef-client -c /etc/chef/client.rb -d -L /var/log/chef/client-stdout.log"
process_name=%(program_name)s
autorestart=false
autostart=false
exitcodes=0,1
stdout_logfile=/var/log/chef/client-stdout.log
stderr_logfile=/var/log/chef/client-stderr.log
