# supervisor/conf/chef-push-jobs-client.ini
[program:chef-push-jobs-client]
command=bash -c "/usr/bin/pushy-client -c /etc/chef/client.rb"
process_name=%(program_name)s
autorestart=false
autostart=false
stdout_logfile=/var/log/chef/push-jobs-client-stdout.log
stderr_logfile=/var/log/chef/push-jobs-client-stderr.log
