# infra-client/conf/supervisor.ini
[program:chef-infra-client]
command=bash -c "chef-client -c /etc/chef/client.rb -d -i 120 -L /var/log/chef/infra-client-stdout.log"
process_name=%(program_name)s
autorestart=false
autostart=false
startsecs=0
startretries=0
exitcodes=0,1
stdout_logfile=/var/log/chef/infra-client-stdout.log
stderr_logfile=/var/log/chef/infra-client-stderr.log
