# ssh/conf/supervisor.ini
[program:ssh]
command=bash -c "sshd -D"
process_name=%(program_name)s
autorestart=false
